# scripts/exit.target
[Unit]
Description=Exit container
AllowIsolate=yes
Requires=poweroff.target
Wants=exit-code.service
